# roles/dnsmasq/files/dnsmasq-cron.service
[Unit]
Description=dnsmasq cron

[Service]
ExecStart=/opt/freifunk/server-scripts/dnsmasq-cron.sh
