# roles/mesh-announce/templates/respondd.service.j2
[Unit]
Description=Respondd

{% for mesh_name in (mesh.keys() |sort()) %}
After=batadv-bat-{{ mesh_name }}.service
After=sys-subsystem-net-devices-vpn\x2d{{ mesh_name }}.device
{% endfor %}

[Service]
ExecStart=/opt/freifunk/ffnord-alfred-announce/respondd.py -d /opt/freifunk/ffnord-alfred-announce/ -b bat-ffv
{%- for domain in (mesh.keys() |sort())  %}
 -i bat-{{ domain }} -i vxlan-{{ domain }} -i vpn-{{ domain }}
{%- endfor %}

Restart=always
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

[Install]
WantedBy=multi-user.target
