# roles/yanic/files/meshviewer.timer
[Unit]
Description=Meshviewer cron Timer

[Timer]
Unit=meshviewer.service
OnBootSec=0
OnUnitActiveSec=30

[Install]
WantedBy=multi-user.target
