# roles/yanic/files/meshviewer.service
[Unit]
Description=Meshviewer cron

[Service]
ExecStart=/opt/freifunk/server-scripts/meshviewer-cron.sh
